# bin/crontab
0 0 * * * curl http://0.0.0.0:5000/smart_locker/gym/reset
