# bin/crontab
0 0 0 ? * * * curl http://0.0.0.0:5000/fitness_locker/reset
